# runtime/ipfs.service
[Unit]
Description=IPFS Daemon

[Service]
User=ubuntu
ExecStart=/usr/local/bin/ipfs daemon
Restart=always

[Install]
WantedBy=multi-user.target